# root/etc/services.d/qwebserver/run
#!/usr/bin/with-contenv bash

cd /app/qwebserver || true

exec \
	s6-setuidgid abc node app.js -c /config/settings.js
