# root/etc/services.d/quassel-web/run
#!/usr/bin/with-contenv bash

cd /app/quassel-web || true

exec \
	s6-setuidgid abc node app.js -c /config/settings-user.js
